# Assignment_4_MP3_Player/overlay/etc/profile.d/prompt.sh
PS1='mp3_shell:\w\$ '
echo "-------------------------------Welcome to Buildroot Mp3 Player-------------------------------"
